# root_overlay/etc/systemd/system/watchdogd.service
[Service]
ExecStart=/usr/sbin/watchdogd -n
