# root_overlay/etc/systemd/system/motion.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/libcamerify /usr/bin/motion